# big-parental-controls/etc/xdg/autostart/big-supervised-indicator.desktop
[Desktop Entry]
Type=Application
Name=Supervised Account Indicator
Name[pt_BR]=Indicador de Conta Supervisionada
Comment=Shows a system tray icon when parental controls are active
Comment[pt_BR]=Mostra um ícone na bandeja do sistema quando o controle parental está ativo
Exec=/usr/bin/big-supervised-indicator
Icon=security-high-symbolic
Terminal=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
X-KDE-autostart-after=panel
